# deploy/systemd/jade-agent-web.service
[Unit]
Description=Jade Agent Web Interface
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=<linux-user>
Group=<linux-user>
WorkingDirectory=/opt/jade_ai_core
EnvironmentFile=/etc/jade-agent/chat.jade-games.com.env
ExecStart=/opt/jade_ai_core/agent_env/bin/python /opt/jade_ai_core/main.py
Restart=always
RestartSec=5
TimeoutStopSec=20
KillSignal=SIGINT
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
